# .cron
15 * * * * cd /app && python update_flagged_dataset.py >> cron.log 2>&1